# customer/update-route53-record.service
[Unit]
After=networking.service

[Service]
ExecStart=/usr/local/bin/update-route53-record.sh

[Install]
WantedBy=default.target
